# config/ansible/roles/ferrarimarco_home_lab_node/templates/monitoring-ont/monitoring-ont-network-config.timer.jinja
[Unit]
Description=Run {{ monitoring_ont_network_config_service_name }} service as scheduled

[Timer]
OnCalendar=*:0/30

[Install]
WantedBy=timers.target
