# docker/ansible/etc/ansible/roles/ferrarimarco_home_lab_node/templates/asuswrt-chkwan/asuswrt-chkwan.timer.jinja
[Unit]
Description=Run ChkWAN service as scheduled

[Timer]
OnCalendar=*:0/1

[Install]
WantedBy=timers.target
