# config/ansible/roles/ferrarimarco_home_lab_node/templates/monitoring/smartmon-prometheus-exporter.timer.jinja
[Unit]
Description=Run smartmon-prometheus-exporter service as scheduled

[Timer]
OnCalendar=*:0/10

[Install]
WantedBy=timers.target
