# docker/ansible/etc/ansible/roles/ferrarimarco_home_lab_node/templates/restic/restic-backup-timer/restic-backup.timer.jinja
[Unit]
Description=Run {{ restic_backup_job_name }} Restic backup job service as scheduled

[Timer]
OnCalendar={{ restic_backup_job_timer_definition }}

[Install]
WantedBy=timers.target
