# config/ansible/roles/ferrarimarco_home_lab_bootstrap/templates/host-configuration/home-lab-host-configuration.timer.jinja
[Unit]
Description=Run the home lab node configuration service as scheduled

[Timer]
OnCalendar={{ home_lab_host_configuration_timer_definition }}

[Install]
WantedBy=timers.target
