# provisioning/common/iot-core/systemd/configuration/lib/systemd/system/start-mqtt-client.service
[Unit]
After=network-online.target
Wants=network-online.target

[Service]
ExecStart=/usr/local/bin/start-mqtt-client.sh

[Install]
WantedBy=multi-user.target
