# config/ansible/roles/ferrarimarco_home_lab_node/templates/monitoring-apt/monitoring-apt.timer.jinja
[Unit]
Description=Run {{ monitoring_apt_tag }} service as scheduled

[Timer]
OnCalendar=*:0/10

[Install]
WantedBy=timers.target
