# docker/ansible/etc/ansible/roles/ferrarimarco_home_lab_node/templates/docker/system-prune-timer/docker-system-prune.timer.jinja
[Unit]
Description=Run docker-system-prune service as scheduled

[Timer]
OnCalendar=Mon *-*-* 05:05:00

[Install]
WantedBy=timers.target
