# docker/ansible/etc/ansible/roles/ferrarimarco_home_lab_node/templates/docker/system-prune-timer/docker-system-prune.service.jinja
[Unit]
Description=Run docker system prune
Requires=docker.service

[Service]
ExecStart=docker system prune --all --force --volumes

[Install]
WantedBy=multi-user.target
